# backend/com.skillglide.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skillglide.backend</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shaansoni/Jobify-1/backend/venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shaansoni/Jobify-1/backend</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/shaansoni/Jobify-1/backend/logs/backend.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shaansoni/Jobify-1/backend/logs/backend_error.log</string>
</dict>
</plist> 